# network/systemd/captive-portal.service
[Unit]
Description=Captive Portal Backend Service
Documentation=https://github.com/captive-portal
After=network.target dnsmasq.service hostapd.service
Wants=dnsmasq.service hostapd.service

[Service]
Type=simple
User=root
Group=root

# Working directory
WorkingDirectory=/opt/captive-portal/backend

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=0.0.0.0
Environment=JWT_SECRET=your-super-secret-jwt-key-change-in-production
Environment=SECURITY_MODE=production

# Node.js executable path (auto-detect WAN interface)
ExecStart=/bin/bash -c "export LAN_INTERFACE=$(ip -o -4 route show to default | awk '{print $5}' | head -n1 || echo eth0); /usr/bin/node src/app.js"

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=captive-portal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/captive-portal/backend/database
ReadWritePaths=/var/log/captive-portal
ReadWritePaths=/var/run/captive-portal

# Resource limits
LimitNOFILE=65535
MemoryMax=512M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
